# level2/hb2023-lv2-reviewer.timer
[Unit]
Description="Red Packet 2023 - Level 2: Auto rejecter"

[Timer]
OnUnitActiveSec=18s
OnBootSec=2min

[Install]
WantedBy=timers.target
